# root/etc/rc.d/S99modemsrv
#!/bin/sh /etc/rc.common
# 自动拨号脚本 — iStoreOS / ImmortalWrt 修复版
# 关键修复：
#   1. 用反引号替代 $() — Ash shell 不支持 $()
#   2. 正确传入 APN 参数 -s cmnet
#   3. 等待 wwan0 就绪后再启动 quectel-CM-M

START=99

start() {
    for i in `seq 1 15`; do
        [ -e /dev/wwan0 ] && break
        sleep 1
    done
    ip link set wwan0 up 2>/dev/null
    quectel-CM-M -s cmnet &
}
