# root/etc/hotplug.d/usb/20-modem_mode
#!/bin/sh
# USB 热插拔脚本 — iStoreOS / ImmortalWrt 修复版
# 插上模组 → 绑定驱动 + 启动拨号
# 拔出模组 → 停止 quectel-CM-M

[ "$ACTION" = add ] || exit 0

# 等待 USB 设备稳定
sleep 2

# 检查 VID:PID
vid=$(cat /sys/bus/usb/devices/$DEVPATH/idVendor 2>/dev/null)
pid=$(cat /sys/bus/usb/devices/$DEVPATH/idProduct 2>/dev/null)

case "$vid:$pid" in
2c7c:0801|2c7c:0900|2c7c:030a)
    logger -t hotplug-modem "Quectel模组已插入 ($vid:$pid)，启动驱动..."

    # 等待 wwan0 出现
    for i in `seq 1 15`; do
        [ -e /dev/wwan0 ] && break
        sleep 1
    done

    # 确保 wwan0 UP
    [ -e /dev/wwan0 ] && ip link set wwan0 up 2>/dev/null

    # 启动 quectel-CM-M 拨号
    [ -e /usr/bin/quectel-CM-M ] && quectel-CM-M -s cmnet &
    ;;
*)
    exit 0
    ;;
esac
